# deploy/systemd/tusommelier-backend.service
[Unit]
Description=TuSommelier FastAPI Backend
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=__SERVER_USER__
Group=__SERVER_USER__
WorkingDirectory=/home/__SERVER_USER__/tusommelier/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/__SERVER_USER__/tusommelier/config/secrets.env
ExecStartPre=/usr/bin/test -x /home/__SERVER_USER__/tusommelier/.venv/bin/python3
ExecStart=/home/__SERVER_USER__/tusommelier/.venv/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
